# xcode/UserData/IB Support/Simulator Devices/FCD16BB7-036F-4201-98F2-2190D314DADE/data/Library/LaunchDaemons/com.apple.CoreSimulator.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.CoreSimulator.bridge</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.CoreSimulator.audio.devices_changed</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.CoreSimulator.audio.devices_changed</string>
			</dict>
			<key>com.apple.backboardd.safetosetidletimer</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.backboardd.safetosetidletimer</string>
			</dict>
			<key>com.apple.springboard.ringerstate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.ringerstate</string>
			</dict>
			<key>com.apple.springboard.volumestate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.volumestate</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.CoreSimulator.CoreSimulatorBridge.lifecycle_support</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.CoreSimulator.host_support</key>
		<true/>
		<key>com.apple.CoreSimulator.pasteboard_promise_support</key>
		<true/>
		<key>com.apple.CoreSimulator.pasteboard_support</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/Library/Developer/PrivateFrameworks/CoreSimulator.framework/Resources/Platforms/iphoneos/usr/libexec/CoreSimulatorBridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/CoreSimulatorBridge</string>
	</array>
	<key>RootedSimulatorPath</key>
	<true/>
</dict>
</plist>
